# KQL/Direct_Autorun_Keys_Modification.yaml
name: Direct Autorun Keys Modification
id: 24357373-078f-44ed-9ac4-6d334a668a11
author: Victor Sergeev, Daniil Yugoslavskiy, oscd.community
date: 2019/10/25
level: medium
description: Detects direct modification of autostart extensibility point (ASEP) in
  registry using reg.exe.
status: test
modified: 2022/08/04
logsource:
  category: process_creation
  product: windows
tactics:
- persistence
relevantTechniques:
- t1547.001
query: 'DeviceProcessEvents

  | where (ProcessCommandLine contains "add" and FolderPath endswith "\\reg.exe")
  and (ProcessCommandLine contains "\\software\\Microsoft\\Windows\\CurrentVersion\\Run"
  or ProcessCommandLine contains "\\software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Userinit"
  or ProcessCommandLine contains "\\software\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell"
  or ProcessCommandLine contains "\\software\\Microsoft\\Windows NT\\CurrentVersion\\Windows"
  or ProcessCommandLine contains "\\software\\Microsoft\\Windows\\CurrentVersion\\Explorer\\User
  Shell Folders" or ProcessCommandLine contains "\\system\\CurrentControlSet\\Control\\SafeBoot\\AlternateShell")'
