# KQL/DLL_Load_via_LSASS.yaml
name: DLL Load via LSASS
id: b3503044-60ce-4bf4-bbcb-e3db98788823
author: Florian Roth (Nextron Systems)
date: 2019-10-16
severity: high
description: Detects a method to load DLL via LSASS process using an undocumented
  Registry key
status: test
modified: 2022-04-21
logsource:
  category: registry_event
  product: windows
tactics:
- Execution
- Persistence
relevantTechniques:
- T1547
query: 'DeviceRegistryEvents

  | where (RegistryKey contains "\\CurrentControlSet\\Services\\NTDS\\DirectoryServiceExtPt"
  or RegistryKey contains "\\CurrentControlSet\\Services\\NTDS\\LsaDbExtPt") and (not(((RegistryValueData
  in~ ("%%systemroot%%\\system32\\ntdsa.dll", "%%systemroot%%\\system32\\lsadb.dll"))
  and InitiatingProcessFolderPath =~ "C:\\Windows\\system32\\lsass.exe")))'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
triggerOperator: GreaterThan
kind: Scheduled
