# KQL/File_Creation_In_Suspicious_Directory_By_Msdt.EXE.yaml
name: File Creation In Suspicious Directory By Msdt.EXE
id: 318557a5-150c-4c8d-b70e-a9910e199857
author: Vadim Varganov, Florian Roth (Nextron Systems)
date: 2022-08-24
severity: high
description: Detects msdt.exe creating files in suspicious directories which could
  be a sign of exploitation of either Follina or Dogwalk vulnerabilities
status: test
modified: 2023-02-23
logsource:
  category: file_event
  product: windows
tactics:
- Persistence
relevantTechniques:
- T1547
query: 'DeviceFileEvents

  | where InitiatingProcessFolderPath endswith "\\msdt.exe" and (FolderPath contains
  "\\Desktop\\" or FolderPath contains "\\Start Menu\\Programs\\Startup\\" or FolderPath
  contains "C:\\PerfLogs\\" or FolderPath contains "C:\\ProgramData\\" or FolderPath
  contains "C:\\Users\\Public\\")'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
triggerOperator: GreaterThan
kind: Scheduled
