# KQL/New_RUN_Key_Pointing_to_Suspicious_Folder.yaml
name: New RUN Key Pointing to Suspicious Folder
id: 02ee49e2-e294-4d0f-9278-f5b3212fc588
author: Florian Roth (Nextron Systems), Markus Neis, Sander Wiebing
date: 2018/08/25
severity: high
description: Detects suspicious new RUN key element pointing to an executable in a
  suspicious folder
status: experimental
modified: 2024/07/16
logsource:
  category: registry_set
  product: windows
tactics:
- Persistence
relevantTechniques:
- T1547
query: 'DeviceRegistryEvents

  | where (((RegistryValueData contains ":\\$Recycle.bin\\" or RegistryValueData contains
  ":\\Temp\\" or RegistryValueData contains ":\\Users\\Default\\" or RegistryValueData
  contains ":\\Users\\Desktop\\" or RegistryValueData contains ":\\Users\\Public\\"
  or RegistryValueData contains ":\\Windows\\Temp\\" or RegistryValueData contains
  "\\AppData\\Local\\Temp\\" or RegistryValueData contains "%temp%\\" or RegistryValueData
  contains "%tmp%\\") or (RegistryValueData startswith "%Public%\\" or RegistryValueData
  startswith "wscript" or RegistryValueData startswith "cscript")) and (RegistryKey
  contains "\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run" or RegistryKey contains
  "\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce")) and (not(((RegistryValueData
  contains "\\AppData\\Local\\Temp\\" or RegistryValueData contains "C:\\Windows\\Temp\\")
  and (RegistryValueData contains "rundll32.exe " and RegistryValueData contains "C:\\WINDOWS\\system32\\advpack.dll,DelNodeRunDLL32")
  and InitiatingProcessFolderPath startswith "C:\\Windows\\SoftwareDistribution\\Download\\"
  and RegistryKey contains "\\Microsoft\\Windows\\CurrentVersion\\RunOnce")))'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
triggerOperator: GreaterThan
kind: Scheduled
