# KQL/Suspicious_Scheduled_Task_Write_to_System32_Tasks.yaml
name: Suspicious Scheduled Task Write to System32 Tasks
id: 80e1f67a-4596-4351-98f5-a9c3efabac95
author: Florian Roth (Nextron Systems)
date: 2021-11-16
severity: high
description: Detects the creation of tasks from processes executed from suspicious
  locations
status: test
modified: 2022-01-12
logsource:
  category: file_event
  product: windows
tactics:
- Execution
- Persistence
relevantTechniques:
- T1053
query: 'DeviceFileEvents

  | where (InitiatingProcessFolderPath contains "\\AppData\\" or InitiatingProcessFolderPath
  contains "C:\\PerfLogs" or InitiatingProcessFolderPath contains "\\Windows\\System32\\config\\systemprofile")
  and FolderPath contains "\\Windows\\System32\\Tasks"'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
triggerOperator: GreaterThan
kind: Scheduled
