# KQL/Suspicious_New_Service_Creation.yaml
name: Suspicious New Service Creation
id: 17a1be64-8d88-40bf-b5ff-a4f7a50ebcc8
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022/07/14
severity: high
description: Detects creation of a new service via "sc" command or the powershell
  "new-service" cmdlet with suspicious binary paths
status: test
modified: 2022/11/18
logsource:
  category: process_creation
  product: windows
tactics:
- persistence
- privilege_escalation
relevantTechniques:
- t1543.003
query: 'DeviceProcessEvents

  | where ((ProcessCommandLine contains "New-Service" and ProcessCommandLine contains
  "-BinaryPathName") or ((ProcessCommandLine contains "create" and ProcessCommandLine
  contains "binPath=") and FolderPath endswith "\\sc.exe")) and (ProcessCommandLine
  contains "powershell" or ProcessCommandLine contains "mshta" or ProcessCommandLine
  contains "wscript" or ProcessCommandLine contains "cscript" or ProcessCommandLine
  contains "svchost" or ProcessCommandLine contains "dllhost" or ProcessCommandLine
  contains "cmd " or ProcessCommandLine contains "cmd.exe /c" or ProcessCommandLine
  contains "cmd.exe /k" or ProcessCommandLine contains "cmd.exe /r" or ProcessCommandLine
  contains "rundll32" or ProcessCommandLine contains "C:\\Users\\Public" or ProcessCommandLine
  contains "\\Downloads\\" or ProcessCommandLine contains "\\Desktop\\" or ProcessCommandLine
  contains "\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\" or ProcessCommandLine
  contains "C:\\Windows\\TEMP\\" or ProcessCommandLine contains "\\AppData\\Local\\Temp")'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
kind: Scheduled
