# KQL/Suspicious_Service_Path_Modification.yaml
name: Suspicious Service Path Modification
id: 138d3531-8793-4f50-a2cd-f291b2863d78
author: Victor Sergeev, oscd.community, Nasreddine Bencherchali (Nextron Systems)
date: 2019/10/21
severity: high
description: Detects service path modification via the "sc" binary to a suspicious
  command or path
status: test
modified: 2022/11/18
logsource:
  category: process_creation
  product: windows
tactics:
- persistence
- privilege_escalation
relevantTechniques:
- t1543.003
query: 'DeviceProcessEvents

  | where (ProcessCommandLine contains "powershell" or ProcessCommandLine contains
  "cmd " or ProcessCommandLine contains "mshta" or ProcessCommandLine contains "wscript"
  or ProcessCommandLine contains "cscript" or ProcessCommandLine contains "rundll32"
  or ProcessCommandLine contains "svchost" or ProcessCommandLine contains "dllhost"
  or ProcessCommandLine contains "cmd.exe /c" or ProcessCommandLine contains "cmd.exe
  /k" or ProcessCommandLine contains "cmd.exe /r" or ProcessCommandLine contains "cmd
  /c" or ProcessCommandLine contains "cmd /k" or ProcessCommandLine contains "cmd
  /r" or ProcessCommandLine contains "C:\\Users\\Public" or ProcessCommandLine contains
  "\\Downloads\\" or ProcessCommandLine contains "\\Desktop\\" or ProcessCommandLine
  contains "\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\" or ProcessCommandLine
  contains "C:\\Windows\\TEMP\\" or ProcessCommandLine contains "\\AppData\\Local\\Temp")
  and (ProcessCommandLine contains "config" and ProcessCommandLine contains "binPath")
  and FolderPath endswith "\\sc.exe"'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
kind: Scheduled
