# KQL/Uncommon_One_Time_Only_Scheduled_Task_At_0000.yaml
name: Uncommon One Time Only Scheduled Task At 00:00
id: 970823b7-273b-460a-8afc-3a6811998529
author: pH-T (Nextron Systems)
date: 2022/07/15
severity: high
description: Detects scheduled task creation events that include suspicious actions,
  and is run once at 00:00
status: test
modified: 2023/02/03
logsource:
  category: process_creation
  product: windows
tactics:
- Execution
- Persistence
- Privilege Escalation
relevantTechniques:
- T1053
query: 'DeviceProcessEvents

  | where (ProcessCommandLine contains "wscript" or ProcessCommandLine contains "vbscript"
  or ProcessCommandLine contains "cscript" or ProcessCommandLine contains "wmic "
  or ProcessCommandLine contains "wmic.exe" or ProcessCommandLine contains "regsvr32.exe"
  or ProcessCommandLine contains "powershell" or ProcessCommandLine contains "\\AppData\\")
  and (FolderPath contains "\\schtasks.exe" or ProcessVersionInfoOriginalFileName
  =~ "schtasks.exe") and (ProcessCommandLine contains "once" and ProcessCommandLine
  contains "00:00")'
eventGroupingSettings:
  aggregationKind: SingleAlert
queryFrequency: P1D
queryPeriod: P1D
enabled: true
entityMappings: null
sentinelEntitiesMappings: null
triggerThreshold: 0
suppressionDuration: PT5H
suppressionEnabled: false
triggerOperator: GreaterThan
kind: Scheduled
